# examples/ai.openclaw.gateway-recovery-watcher.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.gateway-recovery-watcher</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/ABSOLUTE/PATH/TO/recovery/gateway-v04/bin/watch-gateway.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>30</integer>
  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/recovery/gateway-v04</string>
  <key>StandardOutPath</key>
  <string>/ABSOLUTE/PATH/TO/recovery/gateway-v04/logs/watcher.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/ABSOLUTE/PATH/TO/recovery/gateway-v04/logs/watcher.stderr.log</string>
</dict>
</plist>
